# imon-iscsi.service
[Unit]
Description=iMON iSCSI Monitor Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/imon-iscsi.sh
Restart=always
RestartSec=5
User=root
StandardOutput=append:/var/log/imon-iscsi.log
StandardError=append:/var/log/imon-iscsi.log

[Install]
WantedBy=multi-user.target
